# src/main/resources/com/laytonsmith/core/webserver/systemd.service
# Do not edit this file directly! Instead, run `systemctl edit foo` to create an override.

[Unit]
Description=MethodScript WebServer
Documentation=https://methodscript.com/

[Service]
# see man systemd.service 
Type=forking
PIDFile=%%PIDFILE%%
ExecStart=/usr/local/bin/mscript -- webserver --start
ExecReload=/usr/local/bin/mscript -- webserver --recompile
ExecStop=/usr/local/bin/mscript -- webserver --stop
TimeoutStopSec=30
StandardOutput=journal
